# confs/uwsgi.conf
# description "uWSGI server for {{ PROJECT_FILENAME }}"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
    . /etc/environment
    mkdir -p /run/uwsgi 
    chown ubuntu:www-data /run/uwsgi
    /usr/local/bin/uwsgi --ini /etc/uwsgi/sites/elections.app.ini --logto /var/log/elections/uwsgi.log
end script
